# Toxo_AI_code/deploy/iatoxo-backend.service.example
# /etc/systemd/system/iatoxo-backend.service
# Reload: sudo systemctl daemon-reload && sudo systemctl restart iatoxo-backend
# Logs:   sudo journalctl -u iatoxo-backend -f

[Unit]
Description=IA-Toxo FastAPI backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=iatoxo
Group=iatoxo
WorkingDirectory=/srv/iatoxo/backend
EnvironmentFile=/srv/iatoxo/backend/.env
ExecStart=/srv/iatoxo/backend/.venv/bin/uvicorn app.main:app \
          --host 127.0.0.1 \
          --port 8000 \
          --workers 2 \
          --proxy-headers \
          --forwarded-allow-ips 127.0.0.1
Restart=on-failure
RestartSec=3s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/srv/iatoxo/backend

[Install]
WantedBy=multi-user.target
